# conf/supervisor/supervisord.conf
[program:openvpn]
command=openvpn --config /etc/openvpn/server.conf

